# docker/root/etc/services.d/rsyslogd/run
#!/usr/bin/with-contenv bash

exec /usr/sbin/rsyslogd -n -f /config/rsyslog.conf -i /config/rsyslogd.pid
